# contrib/sunstrider/sunrestartloop
### BEGIN INIT INFO
# Provides:          sunLoopRestarter
# Required-Start:    $remote_fs $syslog mysql
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start Suntrider server restarter
# Description:       Start Sunstrider server restarter
### END INIT INFO

#Not a service, this script loops and exec the following commands every 5 seconds IFF maintenance file is set to 0 (or does not exists).
#sunServer start
#sunServer clean_dumps
#Just make sure it's running and starting with the system.

MAINTENANCEFILE=/home/tbc/emu/config/maintenance.conf

while true ; do
  if [ ! -f $MAITENANCEFILE ] || [ $(cat $MAINTENANCEFILE) -eq 0 ] ; then
    service sunserver start &
    service sunserver cleandumps &
  fi
  sleep 5
done


